# rag-system/deploy/qdrant.service
[Unit]
Description=Qdrant Vector Database
After=network.target

[Service]
Type=simple
User=raguser
Group=raguser
WorkingDirectory=/opt/qdrant
ExecStart=/opt/qdrant/qdrant
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
